# profile.d/00-base.sh
# add personal configs
if [ -f "$HOME/.config/$USER/bashrc" ]; then
        source "$HOME/.config/$USER/bashrc"
fi

