# profile.d/99-alias.zsh
#
# 99-alias.zsh
#

alias yolo='yaourt -Syyuua --noconfirm'

